# vidgen.service
[Unit]
Description=Vidgen Video Generation Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vidgen
Environment="PYTHONUNBUFFERED=1"
Environment="HF_HOME=/root/.cache/huggingface"
Environment="PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True"
Environment="TMPDIR=/root/vidgen/tmp"
Environment="TEMP=/root/vidgen/tmp"
Environment="TMP=/root/vidgen/tmp"
ExecStartPre=/bin/mkdir -p /root/vidgen/tmp
ExecStartPre=/bin/chmod 1777 /root/vidgen/tmp
ExecStart=/usr/bin/python3 /root/vidgen/app.py
Restart=always
RestartSec=10
StandardOutput=append:/root/vidgen/vidgen.log
StandardError=append:/root/vidgen/vidgen.log

[Install]
WantedBy=multi-user.target
