# vidgen.service
[Unit]
Description=Vidgen Video Generation Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/termux/vidgen
Environment="PYTHONUNBUFFERED=1"
Environment="HF_HOME=/root/.cache/vidgen_models"
Environment="TRANSFORMERS_CACHE=/root/.cache/vidgen_models"
Environment="PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True"
ExecStartPre=/bin/mkdir -p /root/vidgen/tmp
ExecStartPre=/bin/chmod 1777 /root/vidgen/tmp
ExecStart=/usr/bin/env python3 /root/termux/vidgen/app.py --cache-dir /root/.cache/vidgen_models
Restart=always
RestartSec=10
StandardOutput=append:/root/vidgen/vidgen.log
StandardError=append:/root/vidgen/vidgen.log

[Install]
WantedBy=multi-user.target
