# vidgen.service
[Unit]
Description=Vidgen Video Generation Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vidgen
Environment="PYTHONUNBUFFERED=1"
Environment="HF_HOME=/root/.cache/huggingface"
Environment="TRANSFORMERS_CACHE=/root/.cache/huggingface"
Environment="PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True"
Environment="TOKENIZERS_PARALLELISM=true"
Environment="CUDA_LAUNCH_BLOCKING=0"
Environment="OMP_NUM_THREADS=8"
ExecStart=/usr/bin/python3 /root/vidgen/app.py
Restart=always
RestartSec=10
StandardOutput=append:/root/vidgen/vidgen.log
StandardError=append:/root/vidgen/vidgen.log

[Install]
WantedBy=multi-user.target
